# systemd/.config/systemd/user/clipboard.service
[Unit]
Description=Manage clipboard
After=xorg.target

[Service]
ExecStart=/usr/bin/loliclip
Environment=DISPLAY=:0

[Install]
WantedBy=wm.target
