# systemd/.config/systemd/user/clipboard.service
[Unit]
Description=Manage clipboard

[Service]
ExecStart=/usr/bin/loliclip
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target
