# systemd/.config/systemd/user/graphical.target
[Unit]
Description=Graphical target
Requires=dbus.socket

AllowIsolate=true

[Install]
Alias=default.target
