# systemd/.config/systemd/user/graphical.target
[Unit]
Description=Graphical target
Wants=xorg.target
Requires=dbus.socket
Wants=wm.target

AllowIsolate=true

[Install]
Alias=default.target
