# systemd/.config/systemd/user/mpdris2.service
[Unit]
Description=mpDris2

[Service]
ExecStart=/usr/bin/mpDris2

[Install]
WantedBy=mpd.service
