# systemd/.config/systemd/user/offlineimap.timer
[Unit]
Description=OfflineIMAP Quicksync timer

[Timer]
OnCalendar=*:0/1
Unit=offlineimap.service

[Install]
WantedBy=default.target
