# systemd/.config/systemd/user/printer-applet.service
[Unit]
Description=Printer Applet

[Service]
ExecStart=/usr/bin/system-config-printer-applet

[Install]
WantedBy=graphical.target
