# systemd/.config/systemd/user/sabnzbd.service
[Unit]
Description=Sabnzbd+

[Service]
ExecStart=/usr/bin/sabnzbd -l0

[Install]
WantedBy=default.target
