# systemd/.config/systemd/user/taffybar.service

[Unit]
Description=Taffybar
After=xmonad.service

[Service]
ExecStart=/usr/bin/taffybar

[Install]
WantedBy=xmonad.service
