# systemd/.config/systemd/user/volumeicon.service
[Unit]
Description=Volume Icon

[Service]
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=default.target
