# systemd/.config/systemd/user/wm.target
[Unit]
Description=Window manager
Wants=xorg.target
Requires=dbus.socket

[Install]
